# server/cron/crontab
30 22 15 */3 * /usr/local/bin/python /srv/app/manage.py run_spiders > /proc/1/fd/1 2>&1
